# rootfs/etc/s6-overlay/s6-rc.d/argo-cd-wait/up
#!/usr/bin/env sh

while ! kubectl rollout status --namespace argo-cd deployment argo-cd-server 1> /dev/null 2> /dev/null; do
    sleep 1
done
