# rootfs/etc/s6-overlay/s6-rc.d/argo-cd-refresh/up
#!/usr/bin/env sh

kubectl config set-context --current --namespace=argo-cd

argocd --core app get nginx --refresh
argocd --core app get argo-events --refresh
argocd --core app get argo-rollouts --refresh
argocd --core app get argo-workflows --refresh
