# rootfs/etc/s6-overlay/s6-rc.d/argo-cd-port-forward/run
#!/usr/bin/env sh

exec kubectl \
    port-forward \
    --address="localhost,0.0.0.0" \
    --namespace="argo-cd" \
    service/argo-cd-server \
    80:http
