# rootfs/etc/s6-overlay/s6-rc.d/k3s-copy-kubeconfig/up
#!/usr/bin/env sh

: ${KUBECONFIG}

mkdir -p $(dirname "$KUBECONFIG")

docker cp "k8sage-k3s:${KUBECONFIG}" "${KUBECONFIG}"
