# rootfs/etc/s6-overlay/s6-rc.d/k3s-proxy/run
#!/usr/bin/env sh

: ${K8SAGE_K3S_IP:="$(getent hosts k8sage-k3s | awk '{ print $1 }')"}

exec socat TCP-LISTEN:6443,fork "TCP:${K8SAGE_K3S_IP}:6443"
