# rootfs/etc/s6-overlay/s6-rc.d/git-daemon/run
#!/usr/bin/env sh

exec /usr/bin/git daemon \
    --informative-errors \
    --export-all \
    --base-path="/etc/k8sage/repositories"
