# rootfs/etc/profile.d/helpers
#!/usr/bin/env sh

argocd() {
    command argocd --core "${@:1}"
}
