# rootfs/etc/profile.d/typography.sh
#!/usr/bin/env bash

bold() {
  echo -en "\e[1m"
}

normal() {
  echo -en "\e[0m"
}

red() {
  echo -en "\e[31m"
}

green() {
  echo -en "\e[32m"
}

yellow() {
  echo -en "\e[33m"
}
