# rootfs/etc/s6-overlay/s6-rc.d/k3s-deploy/finish
#!/usr/bin/env sh

docker stop k8sage-k3s
