# rootfs/etc/s6-overlay/s6-rc.d/k3s-wait/up
#!/usr/bin/env sh

while ! docker stats --no-stream k3s 1> /dev/null 2> /dev/null; do
    sleep 1
done
